# deploy/systemd/crossbridge-timeline.service
[Unit]
Description=CrossBridge AI - Function 3 application-timeline API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/crossbridge
ExecStart=/opt/crossbridge/venv/bin/uvicorn server.application_timeline.app:app --host 127.0.0.1 --port 8083
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
