# deploy/systemd/crossbridge-deploy.service
[Unit]
Description=CrossBridge AI - one-button deploy console
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/crossbridge
ExecStart=/opt/crossbridge/venv/bin/uvicorn server.deploy_console.app:app --host 127.0.0.1 --port 8090
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
